# SURFv6revB/project-spec/meta-user/recipes-apps/pueo-squashfs/files/pueo-squashfs.service
[Unit]
Description=pueo-squashfs

[Service]
ExecStart=/usr/bin/pueo-squashfs.sh
StandardOutput=journal

[Install]
WantedBy=multi-user.target
